# deploy/idle-party-rpg.service
[Unit]
Description=Idle Party RPG Game Server
After=network.target

[Service]
Type=simple
User=idlerpg
Group=idlerpg
WorkingDirectory=/opt/idle-party-rpg
ExecStart={{NODE_PATH}} server/dist/index.js
EnvironmentFile=/opt/idle-party-rpg/.env
Environment=NODE_ENV=production
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=idle-party-rpg

[Install]
WantedBy=multi-user.target
